# deploy/backup-db.service
[Unit]
Description=Sarathi-AI — Daily Database & Files Backup
After=network.target

[Service]
Type=oneshot
User=sarathi
Group=sarathi
# Makes BACKUP_GPG_PASSPHRASE + BACKUP_RCLONE_REMOTE (offsite) available to backup.sh
EnvironmentFile=-/opt/sarathi/biz.env
ExecStart=/opt/sarathi/deploy/backup.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sarathi-backup
